# supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/postfix.sh

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:mailman]
command=/mailman.sh
