# supervisord.conf
[supervisord]
nodaemon=true

[program:syslog]
command = /sbin/syslogd -n -O /dev/stdout -S
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
priority = 10

[program:fcgi]
command=/usr/bin/spawn-fcgi -u www-data -g www-data -s /var/run/fcgiwrap.sock -n -- /usr/bin/multiwatch -f 2 -- /usr/sbin/fcgiwrap
priority = 20

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority = 30

[program:postfix]
command=/postfix.sh
priority = 40

[program:mailman]
command=/mailman.sh
priority = 50

[program:cron]
command = /usr/sbin/cron -f -L 15
priority = 60
